# deployments/picoboxd.service
[Unit]
Description=PicoBox Agent Daemon
Documentation=https://github.com/henny-cho/picobox
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Agent manipulates Linux namespaces, cgroups, and mount points; must be root.
User=root
Group=root
ExecStart=/usr/local/bin/picoboxd
Restart=always
RestartSec=3
LimitNOFILE=1048576

# State & logging
StateDirectory=picobox
LogsDirectory=picobox
WorkingDirectory=/var/lib/picobox
Environment=PICOBOX_STORAGE_DIR=/var/lib/picobox

# EnvironmentFile may carry API tokens; leading '-' keeps it optional.
EnvironmentFile=-/etc/default/picoboxd

# Hardening — kept compatible with namespace/mount work the agent performs.
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=true
PrivateTmp=false
ProtectKernelTunables=false
ProtectKernelModules=false
ProtectControlGroups=false
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
RestrictNamespaces=false

[Install]
WantedBy=multi-user.target
